# scripts/systemd/rustynet-relay.service
[Unit]
# D12.d — Rustynet relay sibling service.
#
# Hosts the rustynet-relay binary on Linux for the `relay` and
# `anchor` user-selectable presets in the node-role taxonomy
# (documents/operations/active/NodeRoleTaxonomy_2026-05-21.md).
# Other peers can use this relay as fallback when their direct
# UDP path fails. Forwards ciphertext only; never decrypts peer
# traffic.
#
# This unit is installed alongside rustynetd.service when the
# operator selects a relay-bearing preset (relay or anchor). The
# `rustynet ops install-systemd-relay` helper writes the env file
# at /etc/default/rustynet-relay and enables this unit.
#
# When the operator selects a non-relay-bearing preset (admin,
# client, exit, blind_exit), this unit is disabled + masked by
# `rustynet ops install-systemd-relay --uninstall` so it does not
# auto-start.
Description=Rustynet Relay (ciphertext forwarder for hard-NAT peer pairs)
Documentation=https://github.com/Iwan-Teague/Rustynet
After=network-online.target
Wants=network-online.target
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=simple
EnvironmentFile=-/etc/default/rustynet-relay
# Defaults are conservative: loopback-only bind by default so the
# operator must explicitly widen the bind via /etc/default/rustynet-relay
# before the relay accepts inbound peer-pair sessions from off-host.
Environment=RUSTYNET_RELAY_ID=relay-local
Environment=RUSTYNET_RELAY_BIND=127.0.0.1:4500
Environment=RUSTYNET_RELAY_VERIFIER_KEY=/etc/rustynet/relay-verifier.pub
Environment=RUSTYNET_RELAY_REPLAY_STORE=/var/lib/rustynet/relay.replay
Environment=RUSTYNET_RELAY_PORT_RANGE=40000-49999
Environment=RUSTYNET_RELAY_MAX_SESSIONS_PER_NODE=8
Environment=RUSTYNET_RELAY_MAX_TOTAL_SESSIONS=4096
Environment=RUSTYNET_RELAY_HEALTH_BIND=127.0.0.1:4501
ExecStart=/usr/local/bin/rustynet-relay \
    --relay-id ${RUSTYNET_RELAY_ID} \
    --bind ${RUSTYNET_RELAY_BIND} \
    --verifier-key ${RUSTYNET_RELAY_VERIFIER_KEY} \
    --replay-store ${RUSTYNET_RELAY_REPLAY_STORE} \
    --port-range ${RUSTYNET_RELAY_PORT_RANGE} \
    --max-sessions-per-node ${RUSTYNET_RELAY_MAX_SESSIONS_PER_NODE} \
    --max-total-sessions ${RUSTYNET_RELAY_MAX_TOTAL_SESSIONS} \
    --health-bind ${RUSTYNET_RELAY_HEALTH_BIND}
# Pre-start checks: verifier key and replay store parent must exist.
# The relay refuses to start without them (fail-closed on
# missing trust material).
ExecStartPre=/usr/bin/test -f ${RUSTYNET_RELAY_VERIFIER_KEY}
ExecStartPre=/usr/bin/mkdir -p /var/lib/rustynet
TimeoutStopSec=10s
Restart=on-failure
RestartSec=2s
StateDirectory=rustynet
StateDirectoryMode=0700
RuntimeDirectory=rustynet-relay
RuntimeDirectoryMode=0700
# Hardening — relay forwards ciphertext only, never decrypts. Run
# with the smallest privilege surface that allows UDP bind on the
# configured port.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
SystemCallArchitectures=native
ReadWritePaths=/var/lib/rustynet /run/rustynet-relay
# Bind the configured UDP port even when the relay does not run as
# root by holding CAP_NET_BIND_SERVICE for ports <1024. The default
# bind port (4500) is unprivileged; CAP is retained for operators
# who choose to bind privileged ports.
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
